# archiso-master/archiso/airootfs/etc/systemd/system/autoinstall.service
[Unit]
Description=Arch Linux Auto Installer

[Service]
Type=oneshot
ExecStart=/root/autoinstall.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
